# src/dist/whodis.service
[Unit]
Description=whodis
Documentation=https://github.com/MartinEden/whodis
Wants=network-online.target
After=network-online.target

[Install]
WantedBy=default.target

[Service]
Type=simple
ExecStart=/usr/bin/whodis
